# Svc/Startup/startup.conf
# F Prime startup script
# Launches scripts in /startup

start on filesystem

console output

script
    cd /startup
    ./startup.bash
end script
